# s6-rc.d/ztone/run
#!/command/with-contenv bash
/start_zt1.sh