# s6-rc.d/fileserv/run
#!/command/execlineb -P
/usr/local/bin/gfileserv