# s6-rc.d/fileserv/run
#!/command/with-contenv bash
# make sure we've got it chowned
/bin/chown -R zerotier-one:zerotier-one /opt/key-networks/ztncui
/usr/local/bin/gosu zerotier-one:zerotier-one /usr/local/bin/gfileserv