# s6-rc.d/ztncui/run
#!/command/with-contenv bash
/start_ztncui.sh