# s6-rc.d/entryinit/up
#!/command/with-contenv bash
/start_firsttime_init.sh