# s6-rc.d/entryinit/up
#!/command/execlineb -P
exec /command/with-contenv bash /start_firsttime_init.sh